# equestria-os-save-point/pkg/equestria-os-save-point/usr/lib/systemd/system/equestria-save-point-daily.service
[Unit]
Description=Create daily Equestria Save Point
Wants=local-fs.target
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/lib/equestria-os-save-point/pacman-snapshot.sh D
