# equestria-os-save-point/pkg/equestria-os-save-point/usr/lib/systemd/system/equestria-save-point-daily.timer
[Unit]
Description=Daily Equestria Save Point Backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
