# opt/equestria-os/equestria-os-save-point/pkg/equestria-os-save-point/usr/lib/systemd/system/equestria-auto-restore.service
[Unit]
Description=Equestria OS Auto-Recovery
DefaultDependencies=no
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /usr/lib/equestria-os-save-point/auto-restore.py
ExecStartPost=/usr/bin/systemctl reboot --force
